# crontabs/u0_a128
# Crontab format:
# min hour dayNum month dayOfWeek commands

*/30 * * * * backup
0 12 */2 * * backup --all

# vim: ft=crontab
