# crontabs/u0_a128
# Crontab format:
# min hour dayNum month dayOfWeek commands

# 25 4 * * * play-audio /storage/3039-3361/TheBestMeICanBe/LawOfAttraction/AbundanceAudio1master.mp3

# 35 15 * * * play-audio "/storage/3039-3361/TheBestMeICanBe/Motivation/END LAZINESS - Powerful Motivational Video for Success & Studying 2017.mp4"

0 * * * * backup -s
0 12 */2 * * backup -sa

# vim: ft=crontab
